# docker/supervisor/scheduler.conf
[program:laravel-scheduler]
process_name=%(program_name)s
command=php /var/www/html/artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/scheduler.log
stopwaitsecs=3600

[program:order-regularizer]
process_name=%(program_name)s
command=bash -c 'while true; do php /var/www/html/artisan app:regularize-completed-orders --limit=20 --sleep=3; sleep 300; done'
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/order-regularizer.log
stopwaitsecs=3600
